# deploy/supervisor.example.ini
; kllmgate supervisor 配置示例
;
; 使用方式：
;   1. 复制到 supervisor 配置目录：
;      cp deploy/supervisor.example.ini /etc/supervisor/conf.d/kllmgate.ini
;   2. 根据实际环境修改下方路径和参数
;   3. 加载配置：
;      supervisorctl reread && supervisorctl update

[program:kllmgate]
command=/opt/kllmgate/venv/bin/kllmgate --config /etc/kllmgate/config.toml
directory=/opt/kllmgate
user=kllmgate
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
redirect_stderr=false
stdout_logfile=/var/log/kllmgate/kllmgate.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kllmgate/kllmgate.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; API key 等敏感信息通过环境变量传入
; 多个变量用逗号分隔，不要加引号
environment=
    OPENAI_API_KEY=sk-xxx,
    ANTHROPIC_API_KEY=sk-ant-xxx
